# ovs/supervisord.conf

[program:ovs_monitor]
priority = 4
command = bash -c "/ovs/ovs_monitor.sh"

